# system/etc/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
tty1::respawn:/sbin/getty -nl /sbin/autologin 9600 tty1 xterm-256color
ttyS0::respawn:/sbin/getty -nl /sbin/autologin 9600 ttyS0 xterm-256color
