# system/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:/sbin/getty 115200 tty1 xterm-256color
ttyS0::respawn:/sbin/getty -nl /sbin/autologin 115200 ttyS0 xterm-256color
::ctrlaltdel:/sbin/reboot
